# templates/default/consul-alerts.service.erb
[Unit]
Description=Consul Alerts
Documentation=https://github.com/AcalephStorage/consul-alerts
Requires=network-online.target consul.service
After=network-online.target consul.service

[Service]
User=consul-alert
Group=consul-alert
Environment=GOMAXPROCS=2
EnvironmentFile=-<%= node['consul']['sysconfig_dir'] %>/consul-alerts
ExecStart=<%= node['consul']['bin'] %>/consul-alerts start $OPTIONS
Restart=on-failure
ExecReload=/bin/kill -HUP $MAINPID
KillMode=process

[Install]
WantedBy=multi-user.target
